# templates/systemd/hermes.service
[Unit]
Description=Hermes Agent
Documentation=https://github.com/NousResearch/hermes-agent
Documentation=https://github.com/OnlyTerp/hermes-optimization-guide
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=hermes
Group=hermes
WorkingDirectory=/home/hermes
ExecStart=/usr/local/bin/hermes gateway run
ExecReload=/bin/kill -HUP $MAINPID
Restart=on-failure
RestartSec=5
# `hermes gateway run` exits 1 on a clean SIGTERM shutdown (no dedicated
# success code on stop) — without this, every `systemctl stop`/restart
# shows as a false 'failed' unit state.
SuccessExitStatus=1
StartLimitInterval=300
StartLimitBurst=5

# Environment
EnvironmentFile=-/home/hermes/.hermes/.env
Environment=HOME=/home/hermes
Environment=HERMES_CONFIG=/home/hermes/.hermes/config.yaml
# ProtectHome=read-only below blocks the default ~/.local/state/hermes
# gateway lock dir; redirect it into the one path we do allow writes to
# instead of widening ReadWritePaths.
Environment=XDG_STATE_HOME=/home/hermes/.hermes/xdg_state

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=hermes

# --- Hardening ------------------------------------------------------------
# Filesystem
ProtectSystem=strict
ProtectHome=read-only
# /home/hermes/.ccs holds per-profile CCS state (credentials, session
# history, job queue) needed by the coding-agent-delegate skill's
# CCS-routed Tier 1.
ReadWritePaths=/home/hermes/.hermes /home/hermes/.ccs /tmp
PrivateTmp=true
PrivateDevices=true

# Capabilities
NoNewPrivileges=true
CapabilityBoundingSet=
AmbientCapabilities=

# Namespaces
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectKernelLogs=true
ProtectControlGroups=true
ProtectHostname=true
ProtectClock=true
ProtectProc=invisible
ProcSubset=pid
RestrictNamespaces=true
RestrictRealtime=true
RestrictSUIDSGID=true

# Address families — allow unix for local IPC, inet for HTTPS outbound
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6 AF_NETLINK

# System calls
SystemCallFilter=@system-service
SystemCallFilter=~@privileged @resources @mount @cpu-emulation @debug @reboot @swap
SystemCallArchitectures=native

# Memory / file descriptor limits
MemoryMax=4G
LimitNOFILE=65536
TasksMax=4096

[Install]
WantedBy=multi-user.target
